# systemd/user.service
[Unit]
Description=The Lounge (IRC client)

[Service]
Type=simple
ExecStart=THELOUNGE_HOME=~/.thelounge /usr/bin/thelounge start

[Install]
WantedBy=default.target
