# debian/prerm
#!/bin/bash
deb-systemd-invoke stop thelounge.service
if [ "$1" = "remove" ]; then
	deb-systemd-helper purge thelounge.service
fi
